# chapter7/deploy/circus.conf
description "circusd"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

exec /usr/local/bin/circusd /etc/circus/circusd.ini